# packages/CI_CloudConnector/CI_CloudConnector-0.59.zip/CI_CloudConnector-0.59/CI_CloudConnector/ci_cloudConnectorService
#! /bin/sh
### BEGIN INIT INFO
# Provides:          CI_CloudConnector
# Required-Start:    $all
# Required-Stop:     
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: MCI_CloudConnector
### END INIT INFO

PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/opt/bin
CURRENT_DATE=`date`

. /lib/init/vars.sh
. /lib/lsb/init-functions
# If you need to source some other scripts, do it here

case "$1" in
  start)
    log_begin_msg "Starting CI_CloudConnector service"
    # do something
    cd /
    cd /home/pi/ci_cc
    sudo bash runCloudConnector.sh &

	echo $CURRENT_DATE - Started >> CI_CloudConnectorStartupLog.txt 
    log_end_msg $?
    exit 0
    ;;
  stop)
    log_begin_msg "Stopping CI_CloudConnector"
    # do something to kill the service or cleanup or nothing
	echo $CURRENT_DATE - Stopped >> CI_CloudConnectorStartupLog.txt 

    log_end_msg $?
    exit 0
    ;;
  *)
    echo "Usage: {start|stop}"
    exit 1
    ;;
esac