# packages/melopero-autostart/melopero_autostart-0.0.9.tar.gz/melopero_autostart-0.0.9/melopero-autostart/melopero-autostart.service
[Unit]
Description=melopero-autostart
After=graphical.target

[Service]
Type=forking
ExecStart=/home/melopero-autostart/StartScripts.sh

[Install]
WantedBy=graphical.target
